# protected/data/supervisord.conf
;[unix_http_server]
;file=/tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]
port = 127.0.0.1:9030
username = supervisor
password = foHn83ahull

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
;serverurl=unix:///tmp/supervisor.sock
serverurl = http://localhost:9030
username = supervisor

[group:myhome]
programs=myhome_mail,myhome_image,myhome_sphinx, myhome_service, myhome_export_catalog, myhome_report

[group:test]
programs=test_mail,test_image,test_sphinx, test_service, test_export_catalog, test_report

[program:myhome_mail]
user=myhome
command= /usr/bin/php /var/www/myhome.ru/current/protected/yiic gearmanqueue worker
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/myhome.ru/log/mail.log
stdout_logfile_maxbytes=20MB

[program:myhome_image]
user=myhome
command= /usr/bin/php /var/www/myhome.ru/current/protected/yiic image worker
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/myhome.ru/log/imageGenerator.log
stdout_logfile_maxbytes=20MB

[program:myhome_sphinx]
user=myhome
command= /usr/bin/php /var/www/myhome.ru/current/protected/yiic sphinx worker
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/myhome.ru/log/sphinx.log
stdout_logfile_maxbytes=20MB

[program:myhome_service]
user=myhome
command= /usr/bin/php /var/www/myhome.ru/current/protected/yiic userservice worker
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/myhome.ru/log/userService.log
stdout_logfile_maxbytes=20MB

[program:myhome_export_catalog]
user=myhome
command= /usr/bin/php /var/www/myhome.ru/current/protected/yiic gearmancatalog worker
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/myhome.ru/log/exportCatalog.log
stdout_logfile_maxbytes=20MB

[program:myhome_report]
user=myhome
command= /usr/bin/php /var/www/myhome.ru/current/protected/yiic report worker
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/myhome.ru/log/report.log
stdout_logfile_maxbytes=20MB

;===================================================================

[program:test_sphinx]
user=testmyhome
command= /usr/bin/php /var/www/test.myhome.ru/current/protected/yiic sphinx worker
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/test.myhome.ru/log/sphinx.log
stdout_logfile_maxbytes=20MB

[program:test_image]
user=testmyhome
command= /usr/bin/php /var/www/test.myhome.ru/current/protected/yiic image worker
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/test.myhome.ru/log/imageGenerator.log
stdout_logfile_maxbytes=20MB

[program:test_mail]
user=testmyhome
command= /usr/bin/php /var/www/test.myhome.ru/current/protected/yiic gearmanqueue worker
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/test.myhome.ru/log/mail.log
stdout_logfile_maxbytes=20MB

[program:test_service]
user=testmyhome
command= /usr/bin/php /var/www/test.myhome.ru/current/protected/yiic userservice worker
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/test.myhome.ru/log/userService.log
stdout_logfile_maxbytes=20MB

[program:test_export_catalog]
user=testmyhome
command= /usr/bin/php /var/www/test.myhome.ru/current/protected/yiic gearmancatalog worker
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/test.myhome.ru/log/exportCatalog.log
stdout_logfile_maxbytes=20MB

[program:test_report]
user=testmyhome
command= /usr/bin/php /var/www/test.myhome.ru/current/protected/yiic report worker
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/test.myhome.ru/log/report.log
stdout_logfile_maxbytes=20MB